# .platform/files/profile.d/rbenv.sh
#!/usr/bin/env bash
# rbenv setup
# Make sure we keep RBENV_ROOT up to date with current rbenv
# install directory.
export RBENV_ROOT=/opt/elasticbeanstalk/.rbenv
# export PATH="$RBENV_ROOT/bin:$PATH"
# eval "$(rbenv init -)"
